# Library/LaunchAgents/com.claw.orb-backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claw.orb-backend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/claw2501/services/orb-backend/.venv/bin/uvicorn</string>
        <string>server:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8300</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/claw2501/services/orb-backend</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/orb-backend.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/orb-backend.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DOTENV_PATH</key>
        <string>/Users/claw2501/services/orb-backend/.env</string>
    </dict>
</dict>
</plist>
